# crontab
# crontab
0 7 * * * /app/src/run.sh
# An empty line is required at the end of this file for a valid cron file.